# templates/etc/systemd/system/ipxe-ca-sign.timer.j2
{{ ansible_managed | comment }}

[Unit]
Description=Run the signing script every day at 2 AM

[Timer]
OnCalendar=*-*-* 2:00:00
Persistent=true

[Install]
WantedBy=timers.target
